# templates/varnishncsa.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/local/bin/varnishncsa_systemd_wrapper
